# docker/php_workers/supervisor/supervisor.d/php_workers_amqp.ini
[program:xtags_events_worker]
command=nice -n 15 console messenger:consume "xtags_events" --bus=execute_event.bus --limit=256
process_name=%(program_name)s_%(process_num)02d
numprocs=2
user=www-data
autostart=true
autorestart=unexpected
startsecs=3
startretries=20
stdout_logfile=/usr/src/var/log/%(program_name)s_%(process_num)02d.stdout.log
stdout_logfile_maxbytes=16MB
stdout_capture_maxbytes=1MB
stderr_logfile=/usr/src/var/log/%(program_name)s_%(process_num)02d.stderr.log
stderr_logfile_maxbytes=16MB
stderr_capture_maxbytes=1MB

[program:xtags_commands_low_worker]
command=nice -n 5 console messenger:consume "xtags_commands_low" --bus=execute_command.bus --limit=256
process_name=%(program_name)s_%(process_num)02d
user=www-data
autostart=true
autorestart=unexpected
startsecs=3
startretries=20
stdout_logfile=/usr/src/var/log/%(program_name)s_%(process_num)02d.stdout.log
stdout_logfile_maxbytes=16MB
stdout_capture_maxbytes=1MB
stderr_logfile=/usr/src/var/log/%(program_name)s_%(process_num)02d.stderr.log
stderr_logfile_maxbytes=16MB
stderr_capture_maxbytes=1MB

[program:xtags_commands_high_worker]
command=console messenger:consume "xtags_commands_high" --bus=execute_command.bus --limit=256
process_name=%(program_name)s_%(process_num)02d
user=www-data
autostart=true
autorestart=unexpected
startsecs=3
startretries=20
stdout_logfile=/usr/src/var/log/%(program_name)s_%(process_num)02d.stdout.log
stdout_logfile_maxbytes=16MB
stdout_capture_maxbytes=1MB
stderr_logfile=/usr/src/var/log/%(program_name)s_%(process_num)02d.stderr.log
stderr_logfile_maxbytes=16MB
stderr_capture_maxbytes=1MB
